# OpenRC/dae
#!/sbin/openrc-run
description="dae Service"
command="/usr/local/bin/dae"
command_args="run -c /usr/local/etc/dae/config.dae"
pidfile="/run/${RC_SVCNAME}.pid"
command_background="yes"
output_log="/var/log/dae/access.log"
error_log="/var/log/dae/error.log"
supervisor="supervise-daemon"
rc_ulimit="-n 30000"
rc_cgroup_cleanup="yes"

depend() {
    after docker net net-online sysctl
    need net-online
    use net
}

start_pre() {
    if [ -d /sys/fs/bpf ] && ! mountinfo -q /sys/fs/bpf; then
        eerror "BPF filesystem not mounted, exiting..."
        return 1
    fi
    if [ -d /sys/fs/cgroup ] && ! mountinfo -q /sys/fs/cgroup/; then
        eerror "Control Groups filesystem not mounted, exiting..."
        return 1
    fi
    [ -d "/tmp/dae/" ] || mkdir "/tmp/dae" 
    [ -L "/var/log/dae" ] || ln -s "/tmp/dae" "/var/log/"
    if ! /usr/local/bin/dae validate -c /usr/local/etc/dae/config.dae; then
        eerror "checking config file /usr/local/etc/dae/config.dae failed, exiting..."
        return 1
    fi
}